# esign/esign.conf
start on started framework or framework_ready

env STDERR=/tmp/esign.err

script
	/bin/sh /mnt/us/esign/hide.sh 
end script
          
pre-stop exec /bin/sh /mnt/us/esign/show.sh
          
